# docker/cloudron/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── cron: daily MySQL backup scheduler ───────────────────────────────────────
[program:cron]
command=/usr/sbin/cron -f -L 15
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── nginx: serves the React SPA and proxies /api to Node ─────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── gotenberg: PDF generation via headless Chromium ──────────────────────────
[program:gotenberg]
command=/usr/local/bin/gotenberg --api-port 3001 --chromium-auto-start
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── node: NestJS backend on port 3002 ─────────────────────────────────────────
; All Cloudron env vars are inherited from supervisord (set in start.sh).
; PORT overrides the NestJS listen port so nginx can proxy to 3002.
[program:node]
command=node /app/packages/server/dist/main.js
environment=PORT="3002"
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
